# main.cron
0 2 * * * /usr/bin/python3 /Users/xenahkar/DE-ETL-project/main.py
